# rvc2mqtt/rootfs/etc/services.d/rvc2mqtt/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the rvc2mqtt service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================


# Declare variables
declare interface
declare config_file_path

declare host
declare port
declare username
declare password

# --- LOAD MQTT SETTINGS ---
bashio::log.debug "Setting MQTT details..."
if ! bashio::config.is_empty 'mqtt.host'; then
    host=$(bashio::config 'mqtt.host')
    port=$(bashio::config 'mqtt.port' 1883)
    username=$(bashio::config 'mqtt.username')
    password=$(bashio::config 'mqtt.password')
else
    host=$(bashio::services 'mqtt' 'host')
    port=$(bashio::services 'mqtt' 'port')
    username=$(bashio::services 'mqtt' 'username')
    password=$(bashio::services 'mqtt' 'password')
fi

export "MQTT_HOST=${host}"
export "MQTT_PORT=${port}"
export "MQTT_USERNAME=${username}"
export "MQTT_PASSWORD=${password}"

## Get the values from the user config options.
interface=$(bashio::config 'interface_name')
config_file_path=$(bashio::config 'config_file_path')


bashio::log.info "Starting Up CanBus interface ${interface}"
bashio::log.info "Using config file ${config_file_path}"

# bring up the canbus
exec ip link set ${interface} up
# run the app
exec python3 -m rvc2mqtt.app -c ${config_file_path}
